# run_flask.service
[Unit]
After=network.service

[Service]
ExecStart=/var/www/yohoho/run.sh

[Install]
WantedBy=default.target